# proxy/crontab
# Run cerbot renew non-interactively daily at 2:24.
24 2 * * * certbot renew --quiet